# misc/infrastructure/xonotic-map-indexpages.cron
#!/bin/sh

set -e
. /data/common.sh
LockAndLog ~/"$(basename "${0%.*}").log" 21600
set -x

cd xonotic-map-compiler/misc/infrastructure/xonotic-map-indexpages
git pull
rsync -E build-pages.sh build-index.sh style.css script.js autobuild-bsp-uploader@beta.xonotic.org:
ssh autobuild-bsp-uploader@beta.xonotic.org 'sh build-pages.sh'

